# configs/JupiterNode.service
[Unit]
Description=JupiterSearch text search database node

[Service]
Restart=on-failure
ExecStart=/usr/local/bin/JupiterNode